# 85-snacker-pair.rules
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="10c4", ATTR{idProduct}=="ea60", RUN+="/bin/systemctl start write-stand@stand"


