# 85-snacker-pair.rules
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="1a86", ATTR{idProduct}=="55d4", RUN+="/bin/systemctl start write-stand@stand"

